# LaunchAgents/xyz.aaronstevens.slack-alerts.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  DEPRECATED (2026-05-10). Do not load this plist.

  See scripts/daemon.py header for the full rationale. Short version: a
  launchd-detached daemon does not fire the harness's "task completed"
  notification, so the agent never sees inbound messages even though they
  land in the queue. The working architecture is a foreground listener.sh
  child of the harness with a silent watchdog. See README.md.

  Kept in the repo for forensic value. If you somehow loaded it, unload:

    launchctl unload ~/Library/LaunchAgents/xyz.aaronstevens.slack-alerts.plist
    rm ~/Library/LaunchAgents/xyz.aaronstevens.slack-alerts.plist
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>xyz.aaronstevens.slack-alerts</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/python3</string>
    <string>HOME_PLACEHOLDER/.claude/skills/0-slack-alerts/scripts/daemon.py</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>ThrottleInterval</key>
  <integer>10</integer>

  <key>StandardOutPath</key>
  <string>HOME_PLACEHOLDER/.config/slack-alerts/launchd.out</string>

  <key>StandardErrorPath</key>
  <string>HOME_PLACEHOLDER/.config/slack-alerts/launchd.err</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>HOME</key>
    <string>HOME_PLACEHOLDER</string>
    <key>SLACK_DAEMON_POLL_INTERVAL</key>
    <string>3.0</string>
  </dict>
</dict>
</plist>
